# build/templates/crontab
LOCKDIR=/tmp
* * */4 * * %(user)s flock -w 15 -n $LOCKDIR/sprobot_cron.lock -c "%(virtualenv_path)s/bin/python %(project_root)s/sprobot.py --post"
