# build/templates/crontab
0 */4 * * * %(project_user)s cd %(project_root)s && flock -w 15 -n sprobot_cron.lock -c "%(virtualenv_path)s/bin/python sprobot.py --post"
